# confwatch.path
[Path]
Unit=confwatch.service
PathModified=/etc/entomologist/ento.conf

[Install]
WantedBy=multi-user.target
